# build/openmeteo.slice
[Unit]
Description=Open-Meteo API services

[Slice]
MemoryAccounting=yes
MemoryMax=90%

[Install]
WantedBy=slices.target
