# build/openmeteo-api.service
[Unit]
Description = Run Open-Meteo API

[Service]
EnvironmentFile=/etc/default/openmeteo-api.env
PassEnvironment=VAPOR_ENV API_SYNC_APIKEYS LOG_LEVEL
Type=simple
User=openmeteo-api
Group=openmeteo-api
WorkingDirectory=/var/lib/openmeteo-api/
ExecStart=/usr/local/bin/openmeteo-api serve -b $API_BIND --env $VAPOR_ENV
Restart=always
RestartSec=1
OnFailure=openmeteo-notify@%i.service

[Install]
WantedBy = multi-user.target
