# build/openmeteo-notify@.service
[Unit]
Description=Sent email on failed open-meteo services

[Service]
Type=oneshot
EnvironmentFile=/etc/default/openmeteo-api.env
ExecStart=/usr/local/bin/openmeteo-notify.sh $NOTIFICATION_EMAIL %i
ExecCondition=/bin/test -n "$NOTIFICATION_EMAIL"
User=root
Group=root

[Install]
WantedBy=multi-user.target